# bitnami/apache/rootfs/etc/services.d/apache/finish
#!/usr/bin/with-contenv bash

failcount=0
if [ -f $BITNAMI_APP_DIR/tmp/failcount ]; then
  failcount=$(cat $BITNAMI_APP_DIR/tmp/failcount)
fi

start=$(cat $BITNAMI_APP_DIR/tmp/start)
stop=`date '+%d%H%M%S'`
interval=`expr $stop - $start`
if test $interval -lt 30 ; then
  failcount=`expr $failcount + 1`
else
  failcount=0
fi
echo -n $failcount > $BITNAMI_APP_DIR/tmp/failcount

# bring down container on frequent failures. something is definitely wrong
if test $failcount -ge 3 ; then
  s6-svscanctl -t /var/run/s6/services
fi
