# bitnami/php-fpm/rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bash
set -e
source $BITNAMI_PREFIX/bitnami-utils.sh

mkdir -p $BITNAMI_APP_DIR/tmp
date '+%d%H%M%S' > $BITNAMI_APP_DIR/tmp/start
exec s6-setuidgid $BITNAMI_APP_USER $BITNAMI_APP_DAEMON -F ${EXTRA_OPTIONS:+"$EXTRA_OPTIONS"}
