# bitnami/php-fpm/rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bash
set -e
source $BITNAMI_PREFIX/bitnami-utils.sh

exec s6-setuidgid $BITNAMI_APP_USER php-fpm -F ${EXTRA_OPTIONS:+"$EXTRA_OPTIONS"}
