# containers/mariadb/rootfs/etc/services.d/mariadb/finish
#!/usr/bin/execlineb -S0

if { s6-test $# -ne 0 }
if { s6-test ${1} -eq 255 }
s6-svscanctl -t /var/run/s6/services
