# containers/mariadb/rootfs/etc/services.d/mariadb/run
#!/usr/bin/with-contenv bash
set -e
source $BITNAMI_PREFIX/bitnami-utils.sh

mkdir -p $BITNAMI_APP_DIR/tmp
chown -R $BITNAMI_APP_USER:$BITNAMI_APP_USER $BITNAMI_APP_DIR/tmp || true

if [ -f /tmp/init_mysql.sql ]; then
  PROGRAM_OPTIONS="$PROGRAM_OPTIONS --init-file=/tmp/init_mysql.sql"
fi

exec mysqld.bin $PROGRAM_OPTIONS $EXTRA_OPTIONS
