# containers/nginx/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash
set -e
source $BITNAMI_PREFIX/bitnami-utils.sh
exec nginx -g "daemon off;"
